# auther/morepork_ui.service
[Unit]
Description=Auto auther
Wants=kaiten.service
After=kaiten.service

[Service]
Type=simple
ExecStart=/usr/bin/morepork_ui
StandardOutput=null

[Install]
WantedBy=multi-user.target
